# core/docker/supervisor/supervisord.conf
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:php-fpm]
command=php-fpm -F
priority=10
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=nginx -g "daemon off;"
priority=20
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:cdnlite-scheduler]
command=/bin/sh -c 'if [ "${CDNLITE_SCHEDULER_IDLE:-0}" = "1" ]; then echo "cdnlite scheduler idle"; exec tail -f /dev/null; fi; while true; do php /app/artisan schedule:run || true; sleep "${CDNLITE_SCHEDULER_TICK_SECONDS:-60}"; done'
priority=30
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
